# deploy/systemd/tuning-backend.service
[Unit]
Description=TUNING Backend API Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/tuning/backend
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/path/to/tuning/backend/.env
ExecStart=/usr/bin/node /path/to/tuning/backend/dist/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tuning-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/tuning/backend/data

[Install]
WantedBy=multi-user.target
